# supervisord.conf
[supervisord]
nodaemon=true

[program:samba]
command=/usr/sbin/samba -i

[program:syslog]
command=/usr/sbin/rsyslogd -n
